# bin/smartback4/smartback4-cron.cron
#= /etc/cron.d/smartback4-cron.cron   

0  1 *     *    *    root /var/backups/mysql/smartback4-mysql.sh
0  2 *     *    *    root /root/opensyssetup/bin/smartback4.sh

#

