# bin/smartback4/smartback4-cron.cron
#= /etc/cron.d/smartback3-cron.cron   

0  1 *     *    *    root /var/backups/mysql/smartback3-mysql.sh
0  2 *     *    *    root /usr/local/syssetup/bin/smartback3.sh

#

